# 3.adtac.in/etc/systemd/user/dunst.service
[Unit]
Description=dunst notifications
Wants=graphical.target
After=graphical.target

[Service]
Type=simple
ExecStart=dunst

[Install]
WantedBy=default.target
