# b01.commento.io/etc/systemd/user/psql-backup.service
[Unit]
Description=psql backup service
After=network.service

[Service]
Type=oneshot
ExecStart=%h/bin/psql-backup

[Install]
WantedBy=default.target
